# roles/common/templates/reboot_schedule.timer.j2
[Unit]
Description=Restart daily.

[Timer]
OnCalendar=*-*-* {{reboot_calendar_time}}

[Install]
WantedBy=timers.target
